# package/pi-control-panel_VERSION_ARCHITECTURE/DEBIAN/prerm
#!/bin/sh
# prerm script for pi-control-panel
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

echo 'Executing pre removal steps...'

case "$1" in
    remove|upgrade|deconfigure)
	
		if (systemctl is-active --quiet picontrolpanel) then
			systemctl stop picontrolpanel
		fi
		if (systemctl is-enabled --quiet picontrolpanel) then
			systemctl disable picontrolpanel
		fi
		
		echo -e "\nPiControlPanel was \`$1'd successfully\n"

    ;;
	
    failed-upgrade)
    ;;
	
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
